# desploy/portfolio-neotatto.service
[Unit]
Description=Portfolio Neotatto Flask (Gunicorn)
After=network.target

[Service]
Type=simple
User=gabo
Group=gabo
WorkingDirectory=/home/gabo/portfolio/projects/16-neotatto
EnvironmentFile=-/home/gabo/portfolio/projects/16-neotatto/.env
Environment="PATH=/home/gabo/portfolio/projects/16-neotatto/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStartPre=/usr/bin/mkdir -p /home/gabo/portfolio/projects/16-neotatto/log
ExecStart=/home/gabo/portfolio/projects/16-neotatto/venv/bin/gunicorn \
  --workers 3 \
  --timeout 120 \
  --bind 127.0.0.1:5016 \
  --access-logfile /home/gabo/portfolio/projects/16-neotatto/log/gunicorn-access.log \
  --error-logfile /home/gabo/portfolio/projects/16-neotatto/log/gunicorn-error.log \
  --log-level info \
  app:app

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=portfolio-neotatto
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
